# deploy/meatbash.service
[Unit]
Description=MEATBASH Bun web server
After=network.target

[Service]
Type=simple
User=swader
Group=swader
WorkingDirectory=/var/www/meatbash/code
Environment=PATH=/home/swader/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOST=127.0.0.1
Environment=PORT=3010
ExecStart=/home/swader/.bun/bin/bun run serve
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
